# resources/openWRT/etc/init.d/lise
#!/bin/sh /etc/rc.common

START=51

start() {
	/usr/bin/node /srv/lise/server/index.js
}