# resources/openWRT/etc/init.d/webserial
#!/bin/sh /etc/rc.common

START=99

start() {
	/usr/bin/webserial /dev/ttyATH0 3001 &
}